# rootfs/etc/profile.d/git.sh
#!/bin/ash
# shellcheck shell=ash

alias gc='git clone --depth=1'
