# rootfs/etc/profile.d/ls.sh
#!/bin/ash
# shellcheck shell=ash

alias l='ls --color --group-directories-first -AFhl'
alias la='ls --color --group-directories-first -aFhl'
alias lg='ls --color --group-directories-first -AFhg'
#alias ll='ls --color --group-directories-first -AFho'
#alias lo='ls --color --group-directories-first -AFho'
