# rootfs/etc/profile.d/nano.sh
#!/bin/ash
# shellcheck shell=ash

alias nano='nano -l'
