# rootfs/etc/profile.d/others.sh
#!/bin/ash
# shellcheck shell=ash

alias cdtmp='cd $(mktemp -d)'
