# rootfs/etc/profile.d/tmux.sh
#!/bin/ash
# shellcheck shell=ash

alias tm='tmux'

alias tmls='tmux list-session'

alias tmn='tmux new'
alias tmns='tmux new -s'
alias tmnt='tmux new -t'

alias tma='tmux attach -t'

alias tmk='tmux kill-session'
alias tmkt='tmux kill-session -t'
alias tmks='tmux kill-server'
